# deprecated/scott-filebeat-cron/crontab
* * * * * /bin/date >> /usr/share/filebeat/logs/junk.log
